# config/dns-watch.timer
[Unit]
Description=CaeliCode DNS watcher timer

[Timer]
OnBootSec=5
OnUnitActiveSec=5

[Install]
WantedBy=timers.target
