# config/ssh-bridge.service
[Unit]
Description=CaeliCode SSH agent bridge (Windows → WSL)
After=local-fs.target

[Service]
Type=simple
ExecStart=/opt/caelicode/scripts/ssh-bridge.sh
Restart=on-failure
RestartSec=3
Environment="SSH_AUTH_SOCK=/tmp/caelicode-ssh-agent.sock"

[Install]
WantedBy=default.target
